# config/crontab
*/5 * * * * php /var/www/public/bb-cron.php